# config/udev/91-usb-philips.rules
ACTION=="add", SUBSYSTEM=="usb", ATTRS{product}=="Philips Audio Set", RUN+="/var/lib/soma/usb_pc_link"
