# config/90-dvb-adapter.rules
ACTION=="add", SUBSYSTEM=="usb", ATTRS{product}=="soloHD", RUN+="/bin/chown hts.video /dev/dvb/adapter0"
